# scripts/systemd/rag-search-failure-notify.service
# Desktop notification on rag-search daemon crash
# Install: sudo cp this /etc/systemd/system/ && sudo systemctl daemon-reload
[Unit]
Description=rag-search crash notification
After=rag-search.service

[Service]
Type=oneshot
ExecStart=/usr/bin/notify-send "rag-search crashed" "Daemon stopped unexpectedly. Check: journalctl -u rag-search -n 50"
Environment=DISPLAY=:0
Environment=DBUS_SESSION_BUS_ADDRESS=unix:path=/run/user/1000/bus
