# scripts/systemd/opencode-search.service.d/thermal-max.conf
# Thermal protection drop-in for NVIDIA CUDA GPU (no HW thermal protection)
# Cap GPU temp at 80°C (warn) / 72°C (throttle) via pynvml in the daemon
[Service]
Environment="OPENCODE_THERMAL_MAX_C=80"
Environment="OPENCODE_THERMAL_THROTTLE_C=72"
