# scripts/systemd/rag-search.service.d/cpu-budget.conf
# Kernel-enforced CPU ceiling (HR40) — active work bounded to 1 core; idle/steady-state
# stays < 1% of one core cooperatively (Nice=5/CPUWeight=20) and is proven by the live
# test_cpu_budget.py gate. CPUQuota alone does not imply CPUAccounting (systemd#9647).
# Override CPUQuota= here (e.g. 200%) on hosts where a wider ceiling is acceptable; requires
# the `cpu` controller delegated to the --user manager (see federation-ops-and-invariants.md).
[Service]
CPUAccounting=yes
CPUQuota=100%
